# os-build/backend/go-server/quaiminer-dashboard.service
[Unit]
Description=QuaiMiner CORE OS Dashboard (Go)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=quaiminer
WorkingDirectory=/opt/quaiminer-dashboard
Environment=PORT=3000
Environment=HOST=0.0.0.0
Environment=NODE_RPC_URL=http://localhost:8545
ExecStart=/usr/local/bin/quaiminer-dashboard
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
MemoryMax=100M
CPUQuota=50%

[Install]
WantedBy=multi-user.target

